# root/etc/crontabs/root
*/2 * * * * cd /config/www/ && /usr/bin/python report.py &> /dev/null